# service/service/oled.service
[Unit]
Description=OLED Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/testbed/service
ExecStart=/home/pi/env/bin/python /home/pi/testbed/service/oled_server.py
Restart=always

[Install]
WantedBy=multi-user.target
